# rqdb4ai.service
[Unit]
Description=RQDB4AI API
After=redis-server.service
Requires=redis-server.service

[Service]
Type=simple
WorkingDirectory=%h/work/rqdb4ai
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-%h/work/rqdb4ai/rqdb4ai.env
ExecStart=/usr/bin/python3 -m uvicorn server:app --host 127.0.0.1 --port 18300
Restart=always
RestartSec=5
StandardOutput=append:%h/rqdb4ai.log
StandardError=append:%h/rqdb4ai.log

[Install]
WantedBy=multi-user.target
